# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
pidfile=/tmp/supervisord.pid
logfile=/dev/null
loglevel=info

[program:webdav-mount]
command=/app/scripts/run-webdav-mount.sh
autostart=true
autorestart=true
startsecs=3
startretries=999
priority=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:smbd]
command=/usr/sbin/smbd --foreground --no-process-group --configfile=/etc/samba/smb.conf
autostart=true
autorestart=true
startsecs=3
startretries=999
priority=20
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:nmbd]
command=/usr/sbin/nmbd --foreground --no-process-group --configfile=/etc/samba/smb.conf
autostart=true
autorestart=true
startsecs=3
startretries=999
priority=21
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:backend]
command=/usr/local/bin/python -m micam_sync.main
directory=/app/app
environment=PYTHONPATH="/app/app"
autostart=true
autorestart=true
startsecs=3
startretries=999
priority=30
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
